# .config/systemd/user/excercise.timer
[Unit]
Description=Excercise Motivation

[Timer]
OnCalendar=*-*-* 12:00:00
Unit=excercise.service

[Install]
WantedBy=timers.target
