# .config/systemd/user/excercise.service
[Unit]
Description=Excercise Motivation

[Service]
ExecStart=/usr/bin/notify-send -u critical "Excercise NOW for just 5 minutes!!"
