# .config/systemd/user/late-limit.timer
[Unit]
Description=Late Trigger of poweroff 

[Timer]
OnCalendar=Mon..Sun *-*-* 02:00:00
Unit=limit.service

[Install]
WantedBy=timers.target

