# .config/systemd/user/limit.service
[Unit]
Description=Poweroff system

[Service]
ExecStart=%h/.local/bin/limit
