# .config/systemd/user/notion.service
[Unit]
Description=Remind to fill up Notion daily checklist

[Service]
ExecStart=/usr/bin/notify-send -u critical "Fill up Notion Dailies"
