# .config/systemd/user/study.service
[Unit]
Description=Just start stuyding

[Service]
ExecStart=%h/.local/bin/study
