# .config/systemd/user/update.service
[Unit]
Description=Update system repo packages

[Service]
ExecStart=/usr/bin/sudo pacman -Syu --noconfirm --logfile ~/.local/share/pacman-update.log
